# .docker/s6-rc.d/wait-db/up
#!/command/execlineb -P
foreground { echo Waiting for postgres... }

with-contenv
importas -D "db"    HOST PGHOST
importas -D "5432"  PORT PGPORT

wait-for-it ${HOST}:$PORT --timeout=60 --strict -- echo "postgres is up"
