# .docker/s6-rc.d/metis-backend/run
#!/command/execlineb -P

with-contenv
exec /app/index.py
