# .docker/s6-rc.d/init-metis-backend/up
#!/command/execlineb -P
foreground { echo Initialize metis-backend... }

with-contenv
/app/.docker/init-metis-backend.py
