# .docker/s6-rc.d/aiida/run
#!/command/execlineb -P

with-contenv
verdi daemon start --foreground
