# .docker/s6-rc.d/wait-rmq/up
#!/command/execlineb -P
foreground { echo Wainting for rabbitmq... }

with-contenv
importas -D "rabbitmq"  HOST RMQHOST
importas -D "5672"      PORT RMQPORT

wait-for-it ${HOST}:$PORT --timeout=60 --strict -- echo "rabbitmq is up"
