# .docker/s6-rc.d/init-ssh/up
#!/command/execlineb -P
foreground { echo Initialize ssh... }

with-contenv
/app/.docker/init-ssh.sh
