# .docker/s6-rc.d/yascheduler/run
#!/command/execlineb -P

with-contenv
exec yascheduler -l INFO
