# .docker/s6-rc.d/init-yascheduler/up
#!/command/execlineb -P
foreground { echo Initialize yascheduler... }

with-contenv
/app/.docker/init-yascheduler.py
