# .docker/s6-rc.d/init-aiida/up
#!/command/execlineb -P
foreground { echo Initialize aiida... }

with-contenv
/app/.docker/init-aiida.sh
